# configuration/systemd/sshd.service.d/require_iptables.conf
[Unit]
Requires=iptables.service
Requires=ip6tables.service
